# debian/jabbot.logs-cleanup.cron.d
20 * * * *      nobody find /var/log/jabbot/ -type f -mtime +30 -exec rm {} \;